# config/playlist-vote.service
[Unit]
Description=Playlist Vote

[Service]
ExecStart=/root/app
WorkingDirectory=/root
Restart=always

[Install]
WantedBy=multi-user.target